# env-sandbox/profile.d/bash-aliases-agents
alias aider-local="env OPENAI_API_KEY=sk-empty OPENAI_API_BASE=http://host.docker.internal:8686/v1 aider --model openai/vllm-Qwen3-Coder-30B"
